# main/snapper/rc.snapperd
#!/bin/sh
#
# /etc/rc.d/snapperd: start/stop Snapper service
#

. /etc/rc.subr

NAME="Snapper service"
PROG=/usr/sbin/snapperd

case $1 in
	start)
		msg "Starting $NAME..."
		start_daemon  $PROG
		sleep 1
		;;
	stop)
		msg "Stopping $NAME..."
		stop_daemon $PROG
		;;
	restart)
		$0 stop
		sleep 1
		$0 start
		;;
	status)
		status_daemon $PROG
		;;
	*)
		echo "usage: $0 [start|stop|restart|status]"
		;;
esac

# End of file
