# main/seatd/rc.seatd
#!/bin/sh
#
# /etc/rc.d/seatd: start/stop seatd daemon
#
. /etc/rc.subr

NAME="SEATD daemon"
PROG=/usr/bin/seatd
OPTS="-g seat"

case $1 in
start)
  start_daemon  $PROG  $OPTS &
  ;;
stop)
  stop_daemon  $PROG
  ;;
restart)
  $0 stop
  $0 start
  ;;
status)
  status_daemon  $PROG
  ;;
*)
  echo "usage: $0 [start|stop|restart|status]"
  ;;
esac

# End of file
