# extra/nyx-runit/rc.local
#!/bin/sh
#
# You can add your startup command here
#

# End of file
